# update-motd.d/20-sys-info
#!/bin/bash
printf "\n"
# motd runs as root so we need to use full path
/usr/bin/neofetch --source /home/evan/poop.txt
